# src/rover_web/deploy/systemd/sverh-rover-web-gateway.service
[Unit]
Description=Sverh Rover web gateway
After=network.target sverh-rover-rosbridge.service
Requires=sverh-rover-rosbridge.service

[Service]
Type=simple
WorkingDirectory=%h/sverk_rover
Environment=ROS_AUTOMATIC_DISCOVERY_RANGE=LOCALHOST
# Add Environment=ROS_DOMAIN_ID=<id> if this rover uses a non-default ROS domain.
ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash && source %h/sverk_rover/install/setup.bash && exec ros2 run rover_web web_gateway_node --ros-args --params-file %h/sverk_rover/src/rover_web/config/web.yaml -p identity_file:=%h/sverk_rover/src/rover_web/config/robot_identity.yaml -p rover_config_file:=%h/sverk_rover/src/rover_bringup/config/rover.yaml'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
